# deploy/systemd/kt30-web.service
[Unit]
Description=KT30 Web Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kt30
Group=kt30
WorkingDirectory=/opt/kt30/current
EnvironmentFile=/etc/kt30/kt30.env
Environment=APP_ENV=production
Environment=EMBEDDED_WORKER_ENABLED=false
ExecStart=/opt/kt30/current/.venv/bin/gunicorn --bind 0.0.0.0:8000 --workers 2 --threads 4 --worker-class gthread --timeout 120 --graceful-timeout 120 app:app
Restart=always
RestartSec=5
TimeoutStopSec=30
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
